# openvpn-client-forwarding.service
[Unit]
Description=OpenVPN Client with Port Forwarding
Documentation=man:openvpn(8)
Documentation=https://community.openvpn.net/openvpn/wiki/Openvpn24ManPage
Documentation=https://community.openvpn.net/openvpn/wiki/HOWTO
Wants=network-online.target
After=network-online.target
Before=multi-user.target

[Service]
Type=notify
PrivateTmp=true
WorkingDirectory=/etc/openvpn/client
ExecStartPre=/bin/bash -c 'if [[ ! -f /etc/openvpn/client/openvpn-client-forwarding.env ]]; then echo "Configuration file not found: /etc/openvpn/client/openvpn-client-forwarding.env"; exit 1; fi'
ExecStartPre=/bin/bash -c 'source /etc/openvpn/client/openvpn-client-forwarding.env && if [[ ! -f "$OPENVPN_CONFIG_PATH" ]]; then echo "OpenVPN config file not found: $OPENVPN_CONFIG_PATH"; exit 1; fi'

# Start OpenVPN client
ExecStart=/bin/bash -c 'source /etc/openvpn/client/openvpn-client-forwarding.env && exec /usr/sbin/openvpn --suppress-timestamps --nobind --config "$OPENVPN_CONFIG_PATH"'

# Setup port forwarding after VPN is connected
ExecStartPost=/bin/sleep 5
ExecStartPost=/bin/bash /etc/openvpn/client/setup-port-forwarding.sh /etc/openvpn/client/openvpn-client-forwarding.env setup

# Cleanup port forwarding before stopping VPN
ExecStopPre=/bin/bash /etc/openvpn/client/setup-port-forwarding.sh /etc/openvpn/client/openvpn-client-forwarding.env cleanup

# Stop the OpenVPN process
ExecStop=/bin/kill -TERM $MAINPID
ExecStopPost=/bin/bash /etc/openvpn/client/setup-port-forwarding.sh /etc/openvpn/client/openvpn-client-forwarding.env cleanup

# Restart behavior
Restart=on-failure
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=30

# Security settings
NoNewPrivileges=false
User=root
Group=root

# Capabilities needed for VPN and iptables
CapabilityBoundingSet=CAP_IPC_LOCK CAP_NET_ADMIN CAP_NET_RAW CAP_SETGID CAP_SETUID CAP_SYS_CHROOT CAP_DAC_OVERRIDE CAP_AUDIT_WRITE
AmbientCapabilities=CAP_IPC_LOCK CAP_NET_ADMIN CAP_NET_RAW CAP_SETGID CAP_SETUID CAP_SYS_CHROOT CAP_DAC_OVERRIDE CAP_AUDIT_WRITE

# Device access for TUN interface
DeviceAllow=/dev/net/tun rw
DeviceAllow=/dev/null rw

# Process properties
LimitNPROC=10
LimitNOFILE=1024

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openvpn-client-forwarding

[Install]
WantedBy=multi-user.target